# Pi-Web-Status-Display/networktestservice.conf
description "WAN and Ping testing service"
author "your@email-address.com"

start on runlevel [2345]
stop on runlevel [!2345]

#env

respawn

exec /home/pi/WAN_ping_display.py
