# sensehat/Pi-Temperature-Station/thermoservice.conf
description "My temperature sensor service"
author "your@email-address.com"

start on runlevel [2345]
stop on runlevel [!2345]

env modprobe=/sbin/modprobe

respawn

exec /home/pi/thermometer.py
